# gateway/docker/uvicorn.ini
[program:uvicorn]
command=uvicorn main:app --proxy-headers --host 0.0.0.0 --port 8888
directory=~/.code/sites/python/micro-architecture/gateway/app/
user=pi
stdout_events_enable=true
stderr_events_enable=true
stdout_logfile_maxbytes=1
stderr_logfile_maxbytes=1
stdout_logfiles=/service/logs/app.out.log
stderr_logfiles=/service/logs/app.err.log